# raven_bridge_ble_to_mqtt/DEBIAN/postinst
sudo setcap 'cap_net_raw,cap_net_admin+eip' `which hcitool`
sudo python3 -m pip install pygatt paho-mqtt pexpect retrying kafka-python
sudo systemctl daemon-reload
systemctl start bridgebletomqtt.service
systemctl enable bridgebletomqtt.service
