# raven_bridge_ble_to_mqtt/DEBIAN/prerm
systemctl stop bridgebletomqtt
systemctl disable bridgebletomqtt
